# templates/buildroot/board/copyos/rootfs-overlay/etc/profile.d/copyos.sh
#!/bin/sh
echo "Willkommen bei CopyOS. Verfügbare Befehle: zeig, zu, laufwerke, sicherung-start, hilfe, logbuch"
